# systemd/terranote-monitor-health.timer
[Unit]
Description=Run Terranote Health Monitoring Every 5 Minutes
Documentation=https://github.com/Terranote/terranote-infra
Requires=terranote-monitor-health.service

[Timer]
# Run every 5 minutes
OnCalendar=*:0/5
# If the system was off, run immediately after boot
Persistent=true
# Add some randomization to avoid all systems checking at the same time
RandomizedDelaySec=30

[Install]
WantedBy=timers.target

